# devclaw.service
[Unit]
Description=DevClaw - AI Agent for Tech Teams
After=network.target
Documentation=https://github.com/jholhewres/devclaw

[Service]
Type=simple
User=devclaw
Group=devclaw
WorkingDirectory=/opt/devclaw
ExecStart=/opt/devclaw/devclaw serve --config /etc/devclaw/config.yaml
Restart=on-failure
RestartSec=5

# Variáveis de ambiente (melhor usar EnvironmentFile)
EnvironmentFile=-/etc/devclaw/env

# Segurança (hardening)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/devclaw/sessions /opt/devclaw/data
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

# Limites de recursos
LimitNOFILE=65535
MemoryMax=256M

[Install]
WantedBy=multi-user.target
